# shepherd.service.in
[Unit]
Description=Shepherd LLM Server
Documentation=https://github.com/bigattichouse/shepherd
After=network.target

[Service]
Type=simple
User=shepherd
Group=shepherd
ExecStart=@CMAKE_INSTALL_FULL_BINDIR@/shepherd --server
WorkingDirectory=/var/lib/shepherd
Restart=on-failure
RestartSec=10

# Directories (systemd creates these with correct ownership)
StateDirectory=shepherd
CacheDirectory=shepherd
LogsDirectory=shepherd
ConfigurationDirectory=shepherd

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/lib/shepherd /var/cache/shepherd /var/log/shepherd

[Install]
WantedBy=multi-user.target
